# app-emulation/haspd/files/hasplm.init
#!/sbin/runscript

start() {
	ebegin "Starting HASPLM"
	start-stop-daemon --background --start --pidfile /var/run/hasplm.pid --exec /usr/bin/hasplm
	eend $?
}

stop() {
	ebegin "Stoping HASPLM"
	start-stop-daemon  --stop --pidfile /var/run/hasplm.pid
	eend $?
}
