# app-emulation/haspd/files/winehasp.init
#!/sbin/runscript

start() {
	ebegin "Starting WineHASP"
	start-stop-daemon --background --start --pidfile /var/run/winehasp.pid --exec /usr/bin/winehasp
	eend $?
}

stop() {
	ebegin "Stoping WineHASP"
	start-stop-daemon  --stop --pidfile /var/run/winehasp.pid
	eend $?
}
