# media-gfx/scanbuttond/files/scanbuttond
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License, v2 or later

depend() {
        need modules
}

start() {
        ebegin "Starting scanbuttond"
        start-stop-daemon --start --quiet \
                --exec /usr/bin/scanbuttond
        eend $?
}

stop() {
        ebegin "Stopping scanbuttond"
        start-stop-daemon --stop --quiet --exec /usr/bin/scanbuttond
        eend $?
}
